# vnstat.cron
MAILTO=root
# collect data every 5 minutes
*/5 * * * *  vnstat /usr/sbin/vnstat-cron
# daily reports
5 0 * * * vnstat /usr/bin/vnstat-report
